# release/linux/fpm.before-remove.sh
systemctl stop onekey.service
systemctl disable onekey.service
