# release/linux/fpm.after-install.sh
systemctl enable onekeyd.service
systemctl start onekeyd.service
