# release/linux/fpm.after-install.sh
systemctl enable trezord.service
systemctl start trezord.service
