# release/linux/fpm.before-remove.sh
systemctl stop onekeyd.service
systemctl disable onekeyd.service
