# release/linux/fpm.after-install.sh
systemctl enable onekey.service
systemctl start onekey.service
